# mini-backup.timer
[Unit]
Description=schedule mini-backup

[Timer]
Unit=mini-backup.service
OnCalendar=*-*-* 0:00:00

[Install]
WantedBy=timers.target
